# templates/local-chkbuild.timer.j2
[Unit]
Description=Run chkbuild

[Timer]
{{ chkbuild_timer }}
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
